# overlay/chromeos-base/llm-service/files/llm-service.conf
# llm-service - Offline LLM inference service
#
# Runs llama.cpp server with Qwen2.5 model, exposing HTTP API

description     "LLM inference service with llama.cpp"
author          "chromiumos@example.com"

# Start after system services are ready
start on started system-services
stop on stopping system-services

# Minijail forks off the desired process
expect fork

# Respawn if crashes (but not too aggressively)
respawn
respawn limit 3 60

# OOM score - don't kill this easily (LLM models are expensive to load)
oom score -100

# Run the service with minijail sandboxing
# No bash wrapper - direct execution!
exec minijail0 \
  --config /usr/local/usr/share/minijail/llm-service.minijail.conf \
  -- /usr/local/opt/llm-service/bin/llama-server \
    --model /usr/local/opt/llm-service/models/qwen2.5-1.5b-instruct-q4_k_m.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    --threads 4 \
    --ctx-size 2048 \
    --n-gpu-layers 0 \
    --log-disable \
    --metrics
